# deploy/uaiv-labeler.service
[Unit]
Description=UAIV-Labeler
After=network.target

[Service]
Type=simple
User=YOUR_USER
WorkingDirectory=/opt/UAIV-Labeler
Environment=URBAN_ANNOTATION_HOST=127.0.0.1
Environment=URBAN_ANNOTATION_PORT=7860
Environment=UAIV_BROWSE_ROOTS=/datasets:/opt/UAIV-Labeler/sample_data
Environment=UAIV_DEFAULT_BROWSE_PATH=/opt/UAIV-Labeler/sample_data
Environment=UAIV_QA_ROOT=/opt/UAIV-Labeler/sample_data/qa
ExecStart=/opt/UAIV-Labeler/scripts/start_public.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
